# infra/ansible/slurm/templates/slurmd.service.j2
[Unit]
Description=Slurm node daemon
After=network-online.target munge.service
Wants=network-online.target
Requires=munge.service

[Service]
Type=simple
User=root
Group=root
RuntimeDirectory=slurm
RuntimeDirectoryMode=0755
ExecStart={{ slurm_install_prefix }}/sbin/slurmd -D -s
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
